# overlay/unglue/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

if [ -f /bin/unglue-platform ]; then
    sh /bin/unglue-platform
fi

/bin/unglue-bootstrap

exit 0